# scanner-proxy/systemd/paperdock-scanner-proxy.service
[Unit]
Description=PaperDock Scanner Proxy (FastAPI + SANE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=paperdock
Group=paperdock
WorkingDirectory=/opt/paperdock/scanner-proxy
Environment=PYTHONUNBUFFERED=1
Environment=SCANNER_OUTPUT_DIR=/var/lib/paperdock-scanner/tmp
Environment=SCAN_TIMEOUT_SECONDS=240
Environment=SCAN_FILE_TTL_SECONDS=3600
ExecStart=/opt/paperdock/scanner-proxy/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 3001
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/paperdock-scanner/tmp

[Install]
WantedBy=multi-user.target
